# debian/min.postrm
if which update-rc.d >/dev/null 2>&1; then
    update-rc.d min remove
fi

ldconfig
