# debian/min.prerm
#!/bin/sh
#
# pre-remove script for MIN
#
if which invoke-rc.d >/dev/null 2>&1; then
    invoke-rc.d min stop
else
    /etc/init.d/min stop
fi
true
