# debian/min.postinst
#!/bin/sh
#
# Post-installation script for MIN
#

if which update-rc.d >/dev/null 2>&1; then
    update-rc.d min defaults
fi
if which invoke-rc.d >/dev/null 2>&1; then
    invoke-rc.d min stop
    sleep 1
    invoke-rc.d min start
else
    /etc/init.d/min stop
    sleep 1
    /etc/init.d/min start
fi
ldconfig